# docs/vm-builder-dashboard.service
[Unit]
Description=vm-builder-dashboard
After=network.target

[Service]
Type=simple
User=vm-builder
Group=vm-builder
WorkingDirectory=/opt/vm-builder-dashboard
ExecStart=/opt/vm-builder-dashboard/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8081 --workers 1
Restart=on-failure
RestartSec=5

# Environment overrides — edit as needed
Environment=APISERVER_URL=http://localhost:8080
Environment=DB_PATH=/var/lib/vm-builder-dashboard/db.sqlite3
# Environment=SECRET_KEY=change-me-to-a-random-string

# Harden the service
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/vm-builder-dashboard

[Install]
WantedBy=multi-user.target
